# play.supervisord.conf
[program:play.amitu.com]
command = /home/amitu/play.amitu.com/play.amitu.com
stdout_logfile = /home/amitu/play.amitu.com/supervisor.log
redirect_stderr = true
